# Linux_script/mailhog.service
[Unit]
Description=MailHog
After=network.target

[Service]
ExecStart=/usr/bin/env /usr/local/bin/Mailhog

[Install]
WantedBy=multi-user.target